# deploy/systemd/oddswitch-celery-browser.service
[Unit]
Description=OddSwitch Celery Worker — Browser (Playwright/Bet9ja)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/project/oddswitch
EnvironmentFile=/var/www/project/oddswitch/.env
# Only 1 concurrent browser worker — Playwright is memory-heavy
ExecStart=/var/www/project/oddswitch/venv/bin/celery \
    -A app.workers.celery_app worker \
    -Q browser \
    --concurrency=1 \
    --loglevel=info \
    --logfile=/var/log/oddswitch/celery-browser.log
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oddswitch-celery-browser

[Install]
WantedBy=multi-user.target
