# deploy/systemd/oddswitch-api.service
[Unit]
Description=OddSwitch Translation API (FastAPI/uvicorn)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/project/oddswitch
EnvironmentFile=/var/www/project/oddswitch/.env
ExecStart=/var/www/project/oddswitch/venv/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 8001 \
    --workers 2 \
    --log-level info
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=oddswitch-api

[Install]
WantedBy=multi-user.target
